# internal/systemd/restart.service.tmpl
[Unit]
Description=Periodic restart for Arma Reforger Server - {{.InstanceName}}
After={{.ServiceName}}

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart {{.ServiceName}}
